# s2i/etc/crontab
*/1 * * * * python $APP_ROOT/src/manage.py test