# s2i/etc/crontab
*/2 * * * * python $APP_ROOT/src/manage.py rucio_loader_cron